# conch/conch.service
[Unit]
Description=Conch Shell Phone
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/conch
ExecStart=/usr/bin/python3 /home/pi/conch/conch.py
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
